# systemd/centrunk.cc.service
[Unit]
Description=Centrunk P25 TRS Service (CC)
After=syslog.target network.target

[Service]
User=root
Type=forking
ExecStart=/opt/centrunk/dvmhost/dvmhost -c /opt/centrunk/configs/configCC.yml
Restart=on-abnormal
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target
