# osquery/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": "false",
    "schedule_splay_percent": "10",
    "events_expiry": "3600",
    "database_path": "/var/osquery/osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "enable_monitor": "true",
    "utc": "true"
  },

  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, cpu_physical_cores, physical_memory, hardware_vendor, hardware_model FROM system_info;",
      "interval": 3600,
      "description": "Hardware inventory"
    },
    "os_version": {
      "query": "SELECT name, version, codename, arch FROM os_version;",
      "interval": 3600,
      "description": "OS version details"
    },
    "deb_packages": {
      "query": "SELECT name, version, source FROM deb_packages;",
      "interval": 86400,
      "description": "Installed Debian packages (daily)"
    },
    "users": {
      "query": "SELECT uid, gid, username, shell, directory FROM users WHERE uid >= 1000 OR uid = 0;",
      "interval": 3600,
      "description": "Local user accounts"
    },
    "interface_addresses": {
      "query": "SELECT ia.interface, ia.address, ia.mask, id.mac FROM interface_addresses ia JOIN interface_details id ON ia.interface = id.interface WHERE ia.address != '' AND ia.address NOT LIKE 'fe80%';",
      "interval": 600,
      "description": "Network interface addresses"
    },
    "listening_ports": {
      "query": "SELECT lp.port, lp.protocol, lp.address, p.name, p.path FROM listening_ports lp LEFT JOIN processes p ON lp.pid = p.pid WHERE lp.port != 0;",
      "interval": 300,
      "description": "Security: listening network ports"
    },
    "authorized_keys": {
      "query": "SELECT ak.uid, u.username, ak.key, ak.key_file FROM authorized_keys ak JOIN users u ON ak.uid = u.uid;",
      "interval": 3600,
      "description": "Security: SSH authorized keys"
    },
    "dvmhost_process": {
      "query": "SELECT p.pid, p.name, p.cmdline, p.resident_size, p.user_time, p.system_time, p.start_time FROM processes p WHERE p.name = 'dvmhost';",
      "interval": 60,
      "description": "DVMHost process health"
    },
    "centrunk_services": {
      "query": "SELECT name, sub_state, active_state, load_state FROM systemd_units WHERE name LIKE 'centrunk.%' OR name = 'osqueryd.service' OR name = 'netbird.service';",
      "interval": 120,
      "description": "Centrunk and supporting systemd service states"
    },
    "disk_usage": {
      "query": "SELECT device, path, blocks_size, blocks_available, blocks, inodes, inodes_free FROM mounts WHERE path IN ('/', '/opt', '/var/log');",
      "interval": 300,
      "description": "Disk usage on key mount points"
    },
    "memory_info": {
      "query": "SELECT memory_total, memory_free, memory_available, swap_total, swap_free FROM memory_info;",
      "interval": 120,
      "description": "Memory and swap utilization"
    },
    "cpu_time": {
      "query": "SELECT core, user, nice, system, idle, iowait FROM cpu_time;",
      "interval": 120,
      "description": "Per-core CPU time counters"
    },
    "uptime": {
      "query": "SELECT days, hours, minutes, total_seconds FROM uptime;",
      "interval": 600,
      "description": "System uptime"
    }
  },

  "file_paths": {
    "centrunk_configs": [
      "/opt/centrunk/configs/%%"
    ],
    "centrunk_binaries": [
      "/opt/centrunk/dvmhost/dvmhost"
    ],
    "ssh_authorized_keys": [
      "/home/%/.ssh/authorized_keys",
      "/root/.ssh/authorized_keys"
    ]
  },

  "file_accesses": [
    "centrunk_configs",
    "centrunk_binaries"
  ],

  "packs": {}
}
